# kcptun.ini
[program:kcptun]
command= /root/server_linux_amd64 -c /root/kcptun.json 2>&1
user= root
autostart= true
autorestart= true