# nezha.ini
[program:nezha]
command= /root/nezha-agent -s tz.xiaomingzy.top:443 -p uYxTt3janqjwW7UTjp --tls
user= root
autostart= true
autorestart= true